# deploy/systemd/pipeline.service
[Unit]
Description=Fact Tables Pipeline Runner
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ec2-user
WorkingDirectory=/opt/fact-tables
ExecStart=/opt/fact-tables/scripts/run_pipeline.sh
Restart=on-failure
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
